# debian/inertia_scroller.service
[Unit]
Description=Inertia Scroller Daemon
After=graphical.target

[Service]
Type=simple
ExecStart=/usr/bin/inertia_scroller
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
