# debian/postinst
#!/bin/sh
set -e
if [ "$1" = "configure" ]; then
    systemctl daemon-reload
    systemctl enable momentum_mouse.service
    systemctl start momentum_mouse.service
fi
exit 0
