# debian/prerm
#!/bin/sh
set -e
if [ "$1" = "remove" ]; then
    systemctl stop momentum_mouse.service || true
    systemctl disable momentum_mouse.service || true
    systemctl daemon-reload
fi
exit 0
