# debian/prerm
#!/bin/sh
set -e
if [ "$1" = "remove" ]; then
    systemctl stop inertia_scroller.service || true
    systemctl disable inertia_scroller.service || true
    systemctl daemon-reload
fi
exit 0
